# deploy/lxc/my-onderzoek-pwa.service
[Unit]
Description=My Onderzoek PWA (Next.js)
After=network.target

[Service]
Type=simple
User=nextjs
Group=nextjs
WorkingDirectory=__APP_DIR__
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1
EnvironmentFile=__APP_DIR__/.env
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=__APP_DIR__

[Install]
WantedBy=multi-user.target
